# templates/demo.service.j2
[Unit]
Description=Cleanup demo OCP Projects
After=atomic-openshift-master.service
[Service]
Type=oneshot
ExecStart=/home/{{ openshift_user }}/projects.sh

[Install]
WantedBy=multi-user.target
